# deploy/systemd/mt5-bridge.service
[Unit]
Description=MT5 Bridge Service (TCP + HTTP + MCP)
Documentation=file:///home/ishanp/Documents/GitHub/TradeBridge/deploy/README.md
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ishanp
Group=ishanp
WorkingDirectory=/home/ishanp/Documents/GitHub/TradeBridge

Environment=MT5_TCP_BRIDGE_PORT=8025
Environment=MT5_TCP_BRIDGE_MCP_PORT=8026
Environment=MT5_TCP_BRIDGE_STATUS_PORT=8027
Environment=MT5_GATEWAY_URL=http://127.0.0.1:8020
Environment=MT5_TCP_BRIDGE_ENABLED=true
Environment=PYTHONUNBUFFERED=1

ExecStart=/home/ishanp/Documents/GitHub/TradeBridge/.venv/bin/python -m apps.bridge.main

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt5-bridge

ReadWritePaths=/home/ishanp/Documents/GitHub/TradeBridge
PrivateTmp=true

[Install]
WantedBy=multi-user.target
