# deploy/systemd/mt5-mcp-server.service
[Unit]
Description=TradeBridge MCP Server
Documentation=file:///home/ishanp/Documents/GitHub/TradeBridge/deploy/README.md
After=mt5-tcp-bridge.service mt5-http-gateway.service mt5-terminal.service
Wants=mt5-tcp-bridge.service mt5-http-gateway.service mt5-terminal.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=%h/Documents/GitHub/TradeBridge
Environment=PYTHONPATH=%h/Documents/GitHub/TradeBridge/src
Environment=MT5_GATEWAY_URL=http://127.0.0.1:8020
Environment=MT5_TCP_BRIDGE_ENABLED=true
Environment=MT5_TCP_BRIDGE_PORT=8025
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/Documents/GitHub/TradeBridge/.venv/bin/python -m uvicorn apps.mcp_server.main:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt5-mcp-server

[Install]
WantedBy=default.target
