# deploy/systemd/mt5-terminal.service
[Unit]
Description=MT5 MetaTrader 5 Terminal
Documentation=file:///home/ishanp/Documents/GitHub/TradeBridge/deploy/README.md
StartLimitIntervalSec=300
StartLimitBurst=15

[Service]
Type=simple
Environment=HOME=%h
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=%t
Environment=XAUTHORITY=%h/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
ExecStart=%h/Documents/GitHub/TradeBridge/deploy/systemd/mt5-terminal-wrapper.sh
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt5-terminal

[Install]
WantedBy=default.target
