# deploy/systemd/mt5-tcp-bridge.service
[Unit]
Description=TradeBridge TCP Bridge Server
Documentation=file:///home/ishanp/Documents/GitHub/TradeBridge/deploy/README.md
After=mt5-terminal.service
Wants=mt5-terminal.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=%h/Documents/GitHub/TradeBridge
Environment=PYTHONPATH=%h/Documents/GitHub/TradeBridge/src
Environment=MT5_TCP_BRIDGE_PORT=8025
Environment=MT5_TCP_BRIDGE_HOST=127.0.0.1
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/Documents/GitHub/TradeBridge/.venv/bin/python -m apps.tcp_bridge.main
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt5-tcp-bridge

[Install]
WantedBy=default.target
